# infra/ajyal-bot.service
[Unit]
Description=أجيال المعرفة — Telegram Bot
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/ajyal-bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
User=root

[Install]
WantedBy=multi-user.target
